# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:postgres]
# Note: The path to postgres might change with the postgresql version
command=su - postgres -c "/usr/lib/postgresql/17/bin/postgres -D /var/lib/postgresql/data/pgdata"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:fastapi]
# IMPORTANT: Replace 'app.api.main:app' with the actual path to your FastAPI app instance
command=uvicorn app.main_api:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20

[program:fasthtml]
# IMPORTANT: Replace 'app.ui.main:app' with the actual path to your FastHTML app instance
command=uvicorn app.main_ui:app --host 0.0.0.0 --port 5001
directory=/app
autostart=true
autorestart=true
priority=20

[program:discordbot]
# IMPORTANT: Replace 'app.bot.main' with the actual entrypoint for your bot
command=python app/main_bot.py
directory=/app
autostart=true
autorestart=true
priority=20
